# systemd/myapp-updater.timer
[Unit]
Description=Run Subscription Update Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
